# ansible/roles/vaultwarden/templates/vaultwarden-backup.timer.j2
[Unit]
Description=Vaultwarden backup schedule

[Timer]
OnCalendar={{ vaultwarden_backup_oncalendar }}
Persistent=true
Unit=vaultwarden-backup.service

[Install]
WantedBy=timers.target
