# RALF-lxc-bootstrap-v5.1/scripts/resume.service
[Unit]
Description=RALF Phase 2 Resume
After=network-online.target
Wants=network-online.target
ConditionPathExists=/root/ralf/state/resume-phase2

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ralf-resume.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
