# RALF-lxc-bootstrap-v5.1/scripts/resume.service
[Unit]
Description=RALF Bootstrap Resume v5.1
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ConditionPathExists=/root/.ralf_resume_needed
ExecStart=/root/RALF-lxc-bootstrap-v5.1/scripts/resume.sh

[Install]
WantedBy=multi-user.target
