# ansible/roles/minio/templates/minio.service.j2
[Unit]
Description=MinIO
Documentation=https://min.io/docs/minio/linux/
Wants=network-online.target
After=network-online.target

[Service]
User={{ minio_user }}
Group={{ minio_group }}

EnvironmentFile={{ minio_env_file }}

# systemd does NOT split ${MINIO_OPTS}. Use a shell to expand+split it.
ExecStart=/usr/bin/env bash -lc 'exec {{ minio_bin }} server "${MINIO_VOLUMES}" ${MINIO_OPTS}'

Restart=always
LimitNOFILE=65536
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths={{ minio_data_dir }}

[Install]
WantedBy=multi-user.target
