# ci/gitops-runner.timer
[Unit]
Description=Run RALF GitOps pull every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=gitops-runner.service

[Install]
WantedBy=timers.target
