# ci/gitops-runner.service
[Unit]
Description=RALF GitOps Pull Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=GITOPS_REPO=git@your.git/ralf.git
Environment=GITOPS_PATH=/srv/ralf
Environment=GITOPS_BRANCH=main
ExecStart=/usr/local/bin/gitops-pull.sh

[Install]
WantedBy=multi-user.target
