# ci/ralf-net-scan.timer
[Unit]
Description=Daily RALF network discovery scan

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
